# slyd.conf
description "slyd server"
start on vagrant-mounted
stop on runlevel [!2345]

chdir /vagrant/slyd
script
    PYTHONPATH=/vagrant/slyd:/vagrant/slybot bin/slyd -p 9002 --root ./dist
end script
respawn
