# slyd.conf
description "slyd server"
start on vagrant-mounted
stop on runlevel [!2345]

chdir /vagrant/slyd
exec twistd -n slyd -p 8000
respawn
